# router/templates/sing_box_init.j2
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=99

CONFIG_DIR="{{ shared.sing_box.config_dir }}"
TOKEN="{{ router.sing_box_token }}"
WORKER="https://{{ cloudflare.worker_domain }}"

update_configs() {
{% for f in shared.sing_box.files %}
    wget -4 -q -O "${CONFIG_DIR}/{{ f }}.tmp" "${WORKER}/${TOKEN}/{{ f }}" && \
        mv "${CONFIG_DIR}/{{ f }}.tmp" "${CONFIG_DIR}/{{ f }}" || \
        rm -f "${CONFIG_DIR}/{{ f }}.tmp"
{% endfor %}
}

start_service() {
    nft flush chain ip {{ shared.nftables.table_name_v4 }} output 2>/dev/null
    nft flush chain ip6 {{ shared.nftables.table_name_v6 }} output 2>/dev/null

    update_configs

    nft -f /etc/nftables.d/singbox.nft 2>/dev/null

    procd_open_instance "sing-box_service"
    procd_set_param command {{ shared.sing_box.binary }}
    procd_append_param command run -C "$CONFIG_DIR"
    procd_set_param limits core="unlimited"
    procd_set_param limits nofile="1000000 1000000"
    procd_close_instance
}

reload_service() {
    update_configs

    if {{ shared.sing_box.binary }} check -C "$CONFIG_DIR" 2>/dev/null; then
        kill -HUP $(pgrep -f "sing-box.*run")
    else
        logger -t sing-box "bad config, skipping reload"
    fi
}
